# contrib/automation/launchd/com.notebrain.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notebrain.ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/notebrain</string>
        <string>ingest</string>
    </array>

    <!-- Run every 3 hours (10,800 seconds) -->
    <key>StartInterval</key>
    <integer>10800</integer>

    <!-- Run once when loaded if needed -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Low CPU priority for background execution -->
    <key>Nice</key>
    <integer>19</integer>

    <!-- Redirect logs -->
    <key>StandardOutPath</key>
    <string>/Users/USER/.notebrain/ingest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USER/.notebrain/ingest.log</string>
</dict>
</plist>
